# scripts/com.lunduke.transcripts.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.lunduke.transcripts</string>
    <key>ProgramArguments</key>
    <array>
      <string>/Users/REPLACE_ME/projects/lunduke-transcripts/scripts/run_pipeline.sh</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/REPLACE_ME/projects/lunduke-transcripts</string>
    <key>StartInterval</key>
    <integer>3600</integer>
    <key>StandardOutPath</key>
    <string>/tmp/lunduke-transcripts.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/lunduke-transcripts.err.log</string>
    <key>RunAtLoad</key>
    <true/>
  </dict>
</plist>
